# systemd/proof-countermodel-distill.service
[Unit]
Description=Proof Countermodel Distill Main-Track Pipeline
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/runnel/proof-countermodel-distill
ExecStartPre=/usr/bin/mkdir -p /home/runnel/proof-countermodel-distill/logs/systemd /home/runnel/proof-countermodel-distill/state/systemd/markers
ExecStart=/home/runnel/proof-countermodel-distill/bin/proof_countermodel_service.sh
Environment=PYTHONUNBUFFERED=1
Environment=POCD_GPUS=0,1
Environment=POCD_EVAL_BATCH_SIZE=16
Environment=POCD_EVAL_DATA_SEED=0
Restart=on-failure
RestartSec=60
TimeoutStopSec=900
KillMode=control-group
StandardOutput=append:/home/runnel/proof-countermodel-distill/logs/systemd/pipeline.service.log
StandardError=append:/home/runnel/proof-countermodel-distill/logs/systemd/pipeline.service.log

[Install]
WantedBy=default.target
